# packaging/quikfynd.service
[Unit]
Description=QuikFynd
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/share/quikfynd/QuikFynd/
ExecStart=/usr/share/quikfynd/QuikFynd/QuikFynd
PIDFile=/var/run/quikfynd.pid
Restart=on-abort

[Install]
WantedBy=multi-user.target